# deploy/goodobox.service
[Unit]
Description=Good Open Box - goodobox.com
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/var/www/good-open-box/backend
ExecStart=/usr/bin/node server.js
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
Environment=NODE_ENV=production
Environment=PORT=5000
EnvironmentFile=/var/www/good-open-box/backend/.env

[Install]
WantedBy=multi-user.target
